# packages/arch/cidre-session/fcitx5.service
[Unit]
Description=Fcitx5 Input Method Editor
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/fcitx5 -d
Restart=on-failure

[Install]
WantedBy=graphical-session.target
