# systemd/cidre-desktop-failure.target
[Unit]
Description=Cidre Desktop Failure Target
AllowIsolate=yes
