# downstream/rootfs-overlay/etc/systemd/system/cidre-safe-mode.target
[Unit]
Description=Cidre Safe Mode Target
AllowIsolate=yes
